# Library/LaunchAgents/com.justinpxrk.dark-notify.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.justinpxrk.dark-notify</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/dark-notify</string>
    <string>-c</string>
    <string>{{ .chezmoi.sourceDir }}/Scripts/Themes/on_theme_change.sh</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>KeepAlive</key>
  <true/>
  <key>StandardErrorPath</key>
  <string>{{ .chezmoi.homeDir }}/Library/Logs/dark-notify.log</string>
</dict>
</plist>
