# Library/LaunchAgents/me.justinpxrk.sketchybar-external.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>me.justinpxrk.sketchybar-external</string>
  <key>Program</key>
  <string>/opt/homebrew/bin/sketchybar</string>
  <key>ProgramArguments</key>
  <array>
    <string>external</string>
    <string>-c</string>
    <string>{{ .chezmoi.homeDir }}/.config/sketchybar/sketchybarrc</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>SKETCHYBAR_PROFILE</key>
    <string>external</string>
  </dict>
  <key>KeepAlive</key>
  <true/>
  <key>StandardErrorPath</key>
  <string>{{ .chezmoi.homeDir }}/Library/Logs/me.justinpxrk/sketchybar-external.log</string>
</dict>
</plist>
